# com.rajsingh.ai-news-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rajsingh.ai-news-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>--project</string>
        <string>/Users/rajsingh/Work/ai-news-digest</string>
        <string>python</string>
        <string>-m</string>
        <string>ai_news_digest.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rajsingh/Work/ai-news-digest</string>

    <!-- Run weekdays at 6:50am AEST -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>50</integer>
        </dict>
        <!-- Tuesday -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>50</integer>
        </dict>
        <!-- Wednesday -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>50</integer>
        </dict>
        <!-- Thursday -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>50</integer>
        </dict>
        <!-- Friday -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>6</integer>
            <key>Minute</key>
            <integer>50</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/rajsingh/Work/ai-news-digest/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rajsingh/Work/ai-news-digest/logs/launchd-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
